# Mac-mini/LaunchAgents/com.home-tools.ollama-tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.home-tools.ollama-tracker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/homeserver/Home-Tools/service-monitor/.venv/bin/python3</string>
        <string>/Users/homeserver/Home-Tools/Mac-mini/scripts/ollama_tracker.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/homeserver/Home-Tools/Mac-mini/scripts</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/homeserver/Library/Logs/home-tools/ollama-tracker.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/homeserver/Library/Logs/home-tools/ollama-tracker.log</string>
</dict>
</plist>
